# deploy/google-cloud/kokoro-audiobook-api.service
[Unit]
Description=Kokoro Audiobook API
After=network-online.target
Wants=network-online.target

[Service]
User=comma
Group=comma
WorkingDirectory=/home/comma/kokoro-audiobook
EnvironmentFile=/etc/kokoro-audiobook.env
ExecStart=/home/comma/kokoro-audiobook/.venv/bin/python -m uvicorn audiobook.server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
